# rootfs/etc/services.d/kdf/finish
#!/usr/bin/with-contenv bash
echo "[kdf] Service finished"
