# rootfs/etc/services.d/kdf/run
#!/command/with-contenv bash
# shellcheck shell=bash
set -euo pipefail

echo "[kdf] Starting kdf (KDF)"

# Check if kdf binary exists and is executable
if [[ ! -f /usr/local/bin/kdf ]]; then
    echo "[kdf] ERROR: kdf binary not found at /usr/local/bin/kdf"
    exit 1
fi

if [[ ! -x /usr/local/bin/kdf ]]; then
    echo "[kdf] ERROR: kdf binary is not executable"
    exit 1
fi

# Ensure MM2.json exists by invoking the generator (create/update from /data/options.json)
GEN_PY="/usr/local/bin/generate_mm2.py"
if [[ -f "/opt/kdf-venv/bin/python3" ]]; then
    PYBIN="/opt/kdf-venv/bin/python3"
else
    PYBIN="python3"
fi
if [[ -f "${GEN_PY}" ]]; then
    echo "[kdf] Generating MM2.json via ${GEN_PY}"
    if ! ${PYBIN} "${GEN_PY}"; then
        echo "[kdf] Warning: generate_mm2.py failed; proceeding to check existing MM2.json"
    fi
fi

# Check if MM2.json exists
if [[ ! -f /root/.kdf/MM2.json ]]; then
    echo "[kdf] ERROR: MM2.json not found at /root/.kdf/MM2.json"
    echo "[kdf] Contents of /root/.kdf/:"
    ls -la /root/.kdf/ || echo "Directory does not exist"
    exit 1
fi

echo "[kdf] Binary and config checks passed"

# Start KDF with enhanced logging and persist raw logs to /data/kdf_launch.log
echo "[kdf] Starting KDF with enhanced logging..."
# Ensure log file exists and is writable
mkdir -p /data
touch /data/kdf_launch.log || true
exec bash -c '(/usr/local/bin/kdf 2>&1 | tee -a /data/kdf_launch.log) | while IFS= read -r line; do
    echo "[$(date +"%Y-%m-%d %H:%M:%S")] [kdf] $line"
done'
